# udev/90-usb0-nmcli.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="usb0", RUN+="/usr/bin/nmcli connection up usb0-static"
